# deploy/systemd/ei-api.service
# /etc/systemd/system/ei-api.service
#
# Install:
#   sudo cp deploy/systemd/ei-api.service /etc/systemd/system/ei-api.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now ei-api
#   sudo systemctl status ei-api
#   journalctl -u ei-api -f        # follow logs

[Unit]
Description=EI.Api ASP.NET Core service
After=network.target postgresql.service
Wants=postgresql.service

[Service]
WorkingDirectory=/var/www/ei/server
ExecStart=/usr/bin/dotnet /var/www/ei/server/EI.Api.dll
Restart=always
RestartSec=10
KillSignal=SIGINT
SyslogIdentifier=ei-api
User=www-data
Group=www-data

Environment=ASPNETCORE_ENVIRONMENT=Production
Environment=ASPNETCORE_URLS=http://127.0.0.1:5000
Environment=DOTNET_PRINT_TELEMETRY_MESSAGE=false
Environment=DOTNET_NOLOGO=true

# Hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=full
ProtectHome=true

[Install]
WantedBy=multi-user.target
